# deploy/discord-bot-deploy.service
[Unit]
Description=Auto-deploy Discord Bot (git pull + restart on new commits)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/discord-bot/deploy/deploy.sh
