# deploy/discord-bot.service
[Unit]
Description=Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bot
WorkingDirectory=/opt/discord-bot
ExecStart=/usr/bin/node --env-file=.env index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-bot
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
